# packages/mediacenter/kodi/system.d/hdmi-monitor.service
[Unit]
Description=HDMI connection monitor Kodi
After=kodi.service

[Service]
Type=simple
ExecStart=/usr/bin/hdmi-monitor.start

[Install]
WantedBy=kodi.target
